# src/cputild.service
[Unit]
Description = cputil daemon 

[Service]
ExecStart=/usr/local/bin/cputild

[Install]
WantedBy=default.target
