# src/cputild.service
[Unit]
Description = cputil daemon 

[Service]
ExecStart=/etc/cputild/bin/cputild 

[Install]
WantedBy=default.target
